# update-hosts.service
[Unit]
Description=Update hosts file

[Service]
Type=oneshot
ExecStart=update-hosts
